# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:gophish]
command=/opt/gophish/gophish
directory=/opt/gophish
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
